# vagrant/cookbooks/gearman_dashboard/templates/default/supervisor_program.conf.erb
[program:<%=@name%>]
command=<%=@command%>
directory=<%=@directory%>
user=<%=@user%>
<% if @numprocs > 1 %>
numprocs=<%=@numprocs%>
process_name=%(process_num)s
<% end %>
autostart=true
autorestart=true
stdout_logfile=/var/log/<%=@name%>.log
redirect_stderr=true
stopsignal=QUIT
